# master-cron/crontab
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root

* * * * * root /usr/local/bin/cleanup-agents.sh
